# scripts/leo-clock-wifi.service
[Unit]
Description=Leo Clock — WiFi provisioning check
# Run after NM is up, but before the graphical session that starts the app
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/home/pi/alarm-clock/scripts/wifi-check.sh

[Install]
WantedBy=multi-user.target
